# temperature-backend/src/DEBIAN/postinst
#!/bin/bash
set -e

pip3 install adafruit-circuitpython-dht

systemctl daemon-reload
systemctl enable temperature.service
systemctl start temperature.service

exit 0